# systemd/ori2.service
# ori2 systemd UNIT TEMPLATE. Don't edit in-place.
#
# The bootstrap installer renders it into
#   ~/.config/systemd/user/ori2-<BOT_NAME>.service
# with these placeholders substituted: INSTALL_DIR, BOT_NAME, USER, NPM_PATH,
# PATH (each wrapped in @-signs in the body below).
#
# See INSTALL.md "Headless deployment (Linux, systemd user mode)" for the
# manual substitution + load + restart commands.

[Unit]
Description=Ori2 — @BOT_NAME@ (autonomous assistant)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@INSTALL_DIR@

# ORI2_DAEMON forces daemon mode regardless of stdout TTY status (systemd
# typically already provides no TTY, but be explicit).
Environment=ORI2_DAEMON=true
Environment=NODE_ENV=production
Environment=BOT_NAME=@BOT_NAME@

# PATH is baked at install time so nvm-installed Node resolves. systemd's
# default PATH is /usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin — fine
# for system Node, but invisible to nvm. The bootstrap renders @PATH@ from
# `dirname $(command -v npm)` plus the standard dirs.
Environment=PATH=@PATH@

# PI_CODING_AGENT_DIR pins Pi SDK's global config dir (auth.json, models.json,
# settings.json, themes, debug log) into the bot's own data dir. Without this,
# every bot run by the same OS user shares ~/.pi/agent/ — concurrent writes
# (e.g. pi auth login from one bot while another runs) race. Per-bot makes
# each checkout fully self-contained.
Environment=PI_CODING_AGENT_DIR=@INSTALL_DIR@/data/@BOT_NAME@/.pi-state

# Use the absolute npm path (baked at install time) instead of /usr/bin/env npm
# so nvm/asdf/volta don't break under systemd's minimal env.
ExecStart=@NPM_PATH@ run start

# Restart on crash, give up after 3 consecutive failures within 60s so a
# real bug doesn't crash-loop forever burning CPU.
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Graceful shutdown — give the bot up to 30s to drain Telegram polls and
# stop adapters before forcefully killing.
KillSignal=SIGTERM
TimeoutStopSec=30

# Hardening (user-mode unit; minimal — operator can layer more).
NoNewPrivileges=true
PrivateTmp=true
# data/ is in the install dir, so we don't sandbox $HOME (the bot needs to
# write there). If you move data/ elsewhere, add ReadWritePaths=...

[Install]
WantedBy=default.target
